# mikrotik-backup/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

set -euo pipefail

ssh_pub_key="$(bashio::config ssh_pub_key)"
ssh_pri_key="$(bashio::config ssh_pri_key)"
git_repo="$(bashio::config git_repo)"

mkdir -p ~/.ssh
echo "$ssh_pub_key" > ~/.ssh/id_rsa.pub
echo "$ssh_pri_key" > ~/.ssh/id_rsa
chmod 600 ~/.ssh/id_rsa

if [ ! -f /data/known_hosts ]; then
  echo "github.com ssh-ed25519 AAAAC3NzaC1lZDI1NTE5AAAAIOMqqnkVzrm0SdG6UOoqKLsabgH5C9okWi0dh2l9GKJl" > /data/known_hosts
fi
ln -s /data/known_hosts ~/.ssh/known_hosts

if bashio::config.true "skip_host_key_checking"; then
  bashio::log.warning "skip_host_key_checking is enabled, disable it after first run!"
  echo -e "Host *\n    StrictHostKeyChecking no" > ~/.ssh/config
fi

# Allow deprecated ssh keys
# See https://confluence.atlassian.com/bitbucketserverkb/ssh-rsa-key-rejected-with-message-no-mutual-signature-algorithm-1026057701.html
echo "PubkeyAcceptedKeyTypes +ssh-rsa" >> /etc/ssh/ssh_config

# TODO: if repo exist, check remote url
if [ ! -d /data/repo ]; then
  bashio::log.info "Initial repo cloning..."
  git clone "$git_repo" /data/repo
  mkdir -p /data/repo/mikrotik
fi

echo "SSH known hosts:"
cat ~/.ssh/known_hosts

exec /usr/bin/do_backup
