# mikrotik-backup/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

set -euo pipefail

echo 1
mkdir ~/.ssh
echo 2
echo "$(bashio::config 'ssh_pub_key')"
echo 2.1
bashio::config 'ssh_pub_key'
echo 2.2
bashio::config 'ssh_pub_key' > ~/.ssh/id_rsa.pub
echo 3
bashio::config ssh_pri_key > ~/.ssh/id_rsa
echo 4
chmod 600 ~/.ssh/id_rsa
echo 5

# TODO: if repo exist, check remote url
if [ ! -d /data/repo ]; then
  bashio::log.info "Initial repo cloning..."
  git clone "$(bashio::config git_repo)" /data/repo
fi

echo 6

exec /usr/bin/do_backup
