# mikrotik-backup/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

set -euo pipefail

mkdir ~/.ssh
bashio::config ssh_pub_key > ~/.ssh/id_rsa.pub
bashio::config ssh_pri_key > ~/.ssh/id_rsa
chmod 600 ~/.ssh/id_rsa

# TODO: if repo exist, check remote url
if [ ! -d /data/repo ]; then
  bashio::log.info "Initial repo cloning..."
  git clone "$(bashio::config git_repo)" /data/repo
fi

exec /usr/bin/do_backup
